# kp/knowledge_repo/deploy/kp-knowledge-repo.service
[Unit]
Description=KP Knowledge Repository MCP Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/knowledge_partner
Environment=KP_VIEWER_BASE_URL=https://artifacts.innovatingwithcapella.com
ExecStart=/opt/knowledge_partner/.venv/bin/python -m knowledge_repo
Restart=always
RestartSec=5
StandardOutput=append:/var/log/kp-knowledge-repo/mcp-access.log
StandardError=append:/var/log/kp-knowledge-repo/mcp-error.log

[Install]
WantedBy=multi-user.target
